# frrouting/config/s6/services/pmacctd/run
#!/usr/bin/env bash
exec 2>&1
exec /usr/sbin/pmacctd -f /etc/pmacct/pmacctd.conf
